# debian/nova-api-ec2.postrm
#!/bin/sh -e

if [ "$1" = purge ]; then
	update-rc.d nova-api-ec2 remove
fi

#DEBHELPER#
